# Openfoam_Simulator/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 %(ENV_RESOLUTION)s -ac
autorestart=true
priority=10

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -usepw -create
autorestart=true
priority=20

[program:novnc]
command=/usr/local/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autorestart=true
priority=30

[program:openfoam_simulator]
command=python3 /app/Openfoam_Simulator/src/main.py
environment=DISPLAY=:1,PYTHONPATH=/app/Openfoam_Simulator
directory=/app/Openfoam_Simulator
autorestart=true
priority=40
stdout_logfile=/var/log/openfoam_simulator.log
stderr_logfile=/var/log/openfoam_simulator.err